# scripts/deploy/uu7-backup.timer
[Unit]
Description=Run uu7.io database backup nightly

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
